# deployment/s6-rc.d/ollama/run
#!/command/with-contenv bash

/bin/ollama serve
